# dot_config/systemd/user/protonvpn.service
[Unit]
Description=ProtonVPN
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/protonvpn-app

[Install]
WantedBy=default.target
